# kamiwaza-deb/debian/kamiwaza/DEBIAN/postinst
#!/bin/sh
# postinst script for kamiwaza

set -e

echo "=== Starting Kamiwaza installation ==="
echo "Current directory: $(pwd)"
echo "User: $USER"
echo "SUDO_USER: $SUDO_USER"

# Create installation directory
echo "=== Creating installation directory ==="
INSTALL_DIR="/opt/kamiwaza"
echo "Installation directory: $INSTALL_DIR"
mkdir -p $INSTALL_DIR
cd $INSTALL_DIR
echo "Changed to directory: $(pwd)"

# Download and extract the appropriate tar file
echo "=== Downloading and extracting files ==="
if [ "$(uname)" = "Linux" ]; then
    echo "Linux detected"
    TAR_FILE="kamiwaza-community-0.3.3-UbuntuLinux.tar.gz"
    TAR_URL="https://github.com/kamiwaza-ai/kamiwaza-community-edition/raw/refs/heads/main/kamiwaza-community-0.3.3-UbuntuLinux.tar.gz"
elif [ "$(uname)" = "Darwin" ]; then
    echo "macOS detected"
    TAR_FILE="kamiwaza-community-0.3.3-OSX.tar.gz"
    TAR_URL="https://github.com/kamiwaza-ai/kamiwaza-community-edition/raw/refs/heads/main/kamiwaza-community-0.3.3-OSX.tar.gz"
else
    echo "Unsupported operating system: $(uname)"
    exit 1
fi

echo "Downloading $TAR_FILE from $TAR_URL..."
wget "$TAR_URL"
echo "Extracting $TAR_FILE..."
tar -xvf $TAR_FILE

# Set permissions
echo "=== Setting permissions ==="
chown -R $SUDO_USER:$SUDO_USER $INSTALL_DIR

# Create startup script
echo "=== Creating startup script ==="
cat > /etc/systemd/system/kamiwaza.service << EOF
[Unit]
Description=Kamiwaza Service
After=network.target

[Service]
Type=simple
User=$SUDO_USER
WorkingDirectory=$INSTALL_DIR
ExecStart=/usr/bin/kamiwazad start
Restart=on-failure

[Install]
WantedBy=multi-user.target
EOF

# Find and run the installation script
echo "=== Looking for installation script ==="
cd $INSTALL_DIR
echo "Changed to directory: $(pwd)"
if [ -f "install.sh" ]; then
    echo "Running install.sh with community and license acceptance flags..."
    bash install.sh --community --i-accept-the-kamiwaza-license
else
    echo "Error: install.sh not found in $(pwd)"
    exit 1
fi

# Reload systemd and handle service
echo "=== Setting up service ==="
systemctl daemon-reload
systemctl enable kamiwaza
systemctl start kamiwaza

echo "=== Installation complete! ==="
exit 0 